# bitcoin.service
[Unit]
Description=Bitcoin node container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a bitcoind
ExecStop=/usr/bin/docker stop -t 300 bitcoind

[Install]
WantedBy=multi-user.target
